# deploy/jetson/systemd/vision-stack.service
[Unit]
Description=Vision Stack Services
After=mavproxy.service
Requires=mavproxy.service

[Service]
Type=simple
User=jetson
WorkingDirectory=/home/jetson/drone_system/vision_stack
ExecStart=/usr/bin/python3 -m src.main all --config-dir configs --mode bench_px4_v1_16
Restart=on-failure
RestartSec=5
Environment=GCS_IP=192.168.1.100
Environment=PYTHONPATH=/home/jetson/drone_system/vision_stack

[Install]
WantedBy=multi-user.target
